# config/pico-rgbw-web.service
[Unit]
Description=picow-rgbw web server
After=network.target

[Service]
ExecStart=%h/.local/bin/picow-rgbw-web

[Install]
WantedBy=default.target
